# deploy/supervisor.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=false
minfds=1024
minprocs=200
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── Laravel Queue Worker (default queue) ─────────────────────
[program:restoapp-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php8.3 /var/www/restoapp/backend/artisan queue:work redis \
    --sleep=3 \
    --tries=3 \
    --max-time=3600 \
    --timeout=60
directory=/var/www/restoapp/backend
autostart=true
autorestart=true
startretries=5
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/restoapp-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
user=www-data
stopwaitsecs=60
stopsignal=TERM
killasgroup=true
stopasgroup=true
environment=APP_ENV=production

# ── Laravel Horizon (if using Horizon instead of raw workers) ─
# Uncomment below if you switch to Horizon and comment out the block above
# [program:restoapp-horizon]
# process_name=%(program_name)s
# command=php8.3 /var/www/restoapp/backend/artisan horizon
# directory=/var/www/restoapp/backend
# autostart=true
# autorestart=true
# startretries=5
# numprocs=1
# redirect_stderr=true
# stdout_logfile=/var/log/supervisor/horizon.log
# stdout_logfile_maxbytes=10MB
# stdout_logfile_backups=5
# user=www-data
# stopwaitsecs=3600
# stopsignal=SIGTERM
